# packages/fs-watcher/fs_watcher-1.0.10-py2.py3-none-any.whl/fs_watcher-1.0.10.data/data/share/fs-watcher/init/fs-watcher.service
[Unit]
Description=filesystem watcher
After=basic.target

[Service]
ExecStart=/usr/sbin/watcher -c /etc/watcher.ini start
ExecStop=/usr/sbin/watcher -c /etc/watcher.ini stop
Type=simple
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
